# deploy/systemd/bv-inbox-sync.timer
[Unit]
Description=Periodically sync encrypted-inbox addresses for Postfix
[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
[Install]
WantedBy=timers.target
